# ubos-networking/systemd/ubos-networking-gateway.service
[Unit]
Description=UBOS networking (gateway)
Requires=\
    avahi-daemon.service \
    network.target \
    systemd-networkd.service \
    systemd-resolved.service \
    systemd-timedated.service \
    ubos-admin.service \
    ubos-nftables.service
Before=\
    avahi-daemon.service \
    systemd-networkd.service \
    systemd-resolved.service \
    systemd-timedated.service \
    ubos-admin.service \
    ubos-nftables.service
Conflicts=\
    cloud-init.service \
    ubos-networking-client.service \
    ubos-networking-cloud.service \
    ubos-networking-off.service \
    ubos-networking-standalone.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStartPre=/usr/share/ubos-admin/bin/ubos-service-prestart ubos-networking-gateway
ExecStart=/usr/share/ubos-networking/bin/activate-netconfig gateway
ExecStopPost=/usr/share/ubos-admin/bin/ubos-service-poststop ubos-networking-gateway

[Install]
WantedBy=multi-user.target
