# ubos-networking/systemd/ubos-networking-off.service
[Unit]
Description=UBOS networking (off)
Requires=\
    ubos-admin.service
Conflicts=\
    avahi-daemon.service \
    cloud-init.service \
    systemd-networkd.service \
    systemd-resolved.service \
    systemd-timedated.service \
    ubos-networking-client.service \
    ubos-networking-cloud.service \
    ubos-networking-gateway.service \
    ubos-networking-standalone.service \
    ubos-nftables.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStartPre=/usr/share/ubos-admin/bin/ubos-service-prestart ubos-networking-off
ExecStart=/usr/share/ubos-networking/bin/activate-netconfig off
ExecStopPost=/usr/share/ubos-admin/bin/ubos-service-poststop ubos-networking-off

[Install]
WantedBy=multi-user.target
