# ubos-admin/systemd/ubos-admin.service
[Unit]
Description=Configure or reconfigure ubos-admin on startup
Wants=network.target cronie.service ntpd.service
Requires=ubos-mysqld.service ubos-httpd.service
# ubos-admin wants to set a root password for mysql, so start mysqld first
After=ubos-mysqld.service
# ubos-admin wants to change the Apache config, so start Apache later
Before=ubos-httpd.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/share/ubos-admin/bin/initialize

[Install]
WantedBy=multi-user.target
